# systemd/yqgl-asr.service
[Unit]
Description=YQGL Qwen3-ASR FastAPI server
After=network.target

[Service]
User={{USER}}
WorkingDirectory=/srv/yqgl/asr_service
Environment=QWEN_ASR_MODEL_PATH={{HOME}}/.cache/modelscope/hub/models/Qwen/Qwen3-ASR-1___7B
Environment=QWEN_ASR_DEVICE=cuda:0
Environment=QWEN_ASR_DTYPE=bfloat16
Environment=PYTHONUSERBASE={{HOME}}/.local
ExecStart={{PY313}} -m uvicorn server:app --host 127.0.0.1 --port 8001 --workers 1 --timeout-keep-alive 120
Restart=on-failure
RestartSec=15
StandardOutput=append:/srv/yqgl/data/asr.log
StandardError=append:/srv/yqgl/data/asr.log

[Install]
WantedBy=multi-user.target
